# data/nocc-daemon.socket
[Unit]
Description=nocc-daemon socket

[Socket]
ListenStream=/run/nocc-daemon.sock

[Install]
WantedBy=sockets.target
