# data/nocc-server.service
[Unit]
Description=Nocc: A Distributed Compilation Server
After=network.target

[Service]
ExecStart=/usr/bin/nocc-server

[Install]
WantedBy=multi-user.target
