# data/nocc-daemon.socket
[Unit]
Description=nocc-daemon socket
Requires=nocc-daemon.service

[Socket]
ListenStream=/run/nocc-daemon.sock

[Install]
WantedBy=sockets.target
